# etc/gunicorn/supervisor-gunicorn_djskel.conf
[program:djskel]
command = /home/webapps/Projects/djskel/bin/gunicorn.sh                       ; Command to start app
user = djskel                                                                        ; User to run as
stdout_logfile = /home/webapps/Projects/djskel/var/log/gunicorn-supervisor.log      ; Where to write log messages
redirect_stderr = true                                                                      ; Save stderr in the same log
